# build/logio/supervisor.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=1MB
loglevel=warn
user=logio
logfile=/home/logio/supervisord/supervisord.log
pidfile=/home/logio/supervisord/supervisord.pid

loglevel=error

[program:logio]
command=log.io-server
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_maxbytes=0
stderr_logfile_maxbytes = 0
stdout_logfile_maxbytes = 0
startsecs=0
autorestart=true
startretries=1

[program:harvester]
command=log.io-harvester
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
startsecs=0
autorestart=true
startretries=1